# bitlbee/program-bitlbee.conf
[program:bitlbee]
command = /usr/local/bin/launch-bitlbee.sh
