# znc/src/program-znc.conf
[program:znc]
user = znc
command = znc -d /var/lib/znc -f -n
